# openbts_systemd_scripts-master/systemd/smqueue.service
[Unit]
Description=smqueue

[Service]
Type=simple
ExecStart=/OpenBTS/smqueue

[Install]
WantedBy=multi-user.target
